# docker/crono-supervisord.conf
[supervisord]
nodaemon=true

[program:crono]
command=/bin/bash -c "CATALOG_CRON=10.seconds bundle exec crono"
