# docker/crono-supervisord.conf
[supervisord]
nodaemon=true

[program:crono]
command=/bin/bash -c "bundle exec crono"
